# othersrc/usr-lib/systemd/system/user-runtime-dir@.service.d/override.conf
[Unit]
ConditionPathIsMountPoint=!/run/user/%i
