# othersrc/usr-lib/systemd/system/user-runtime-dir@.service.d/override.conf
[Service]
[Service]
ExecStart=
ExecStart=/lib/systemd/systemd-user-runtime-dir start %i
ExecStart=/usr/lib/bottle-imp/map-user-runtime-dir.sh %i
ExecStop=
ExecStop=/usr/lib/bottle-imp/unmap-user-runtime-dir.sh %i
ExecStop=/lib/systemd/systemd-user-runtime-dir stop %i

