# othersrc/usr-lib/systemd/system/imp-remount-root-shared.service
[Unit]
Description=bottle-imp - Remount Root Filesystem Shared
DefaultDependencies=no
Before=sysinit.target

[Service]
Type=oneshot
ExecStart=/bin/mount --make-rshared /

[Install]
WantedBy=sysinit.target
